# broadlink_ac_mqtt/root/etc/services.d/ac2mqtt/run
#!/command/with-contenv bashio

exec \
  python3 /app/ac2mqtt/monitor.py \
  -ms "$(bashio::services mqtt host && true)" \
  -mp "$(bashio::services mqtt port && true)" \
  -mU "$(bashio::services mqtt username && true)" \
  -mP "$(bashio::services mqtt password && true)" \
  -c /config/broadlink_ac_mqtt.yaml
